# config/messenger-worker.ini
[program:messenger-consume]
command=php /var/www/froqdam.com/bin/console messenger:consume pimcore_core pimcore_maintenance --memory-limit=250M --time-limit=3600
user=ubuntu
numprocs=2
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stderr_logfile=/var/log/supervisord_froqdam.com.err.log

[program:maintenance]
command=bash -c 'sleep 3600 && exec php /var/www/html/froqdam.com/bin/console pimcore:maintenance --async'
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
